# config/etc/services.d/workerman/run
#!/bin/sh

cd /var/www/html

php start.php start